# bot.conf
[program:server]
directory = /home/ubuntu/bot
command =/home/ubuntu/bot/server
autostart = true
autorestart = true
startsecs = 5
startretries = 30
user = nobody
